# lib/Pushbullet-BTN.conf
description "Pushbullet-BTN-Integration"
author      "ther1v3r"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 99 5

script
    exec /usr/local/bin/node /opt/Pushbullet-BTN/lib/app.js >> /var/log/node.log 2>&1
end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
   # /root/bin/hoptoad.sh "Pushbullet-BTN has started!"
end script
